# templates/101_usb_removed.rules
SUBSYSTEM=="usb"
ACTION=="remove"
ENV{DEVTYPE}=="usb_device"
ENV{PRODUCT}="{{ item.env_product }}"
RUN+="{{ item.path_script }}"
